# docker/etc/s6/smartdns/run
#!/bin/bash

# 配置
source ./setup

# 启动服务器
su-exec "${USERNAME}" smartdns -c ${CONFIG_FILEPATH} -p - -f
